# wififix/etc/systemd/system/wififix.service
[Unit]
Description=Fix wifi issues
After=network.target suspend.target

[Service]
User=root
Type=idle
ExecStart=/usr/local/bin/wififix

[Install]
WantedBy=network.target suspend.target
